# dashboards/supervisord.conf
; Copyright (c) 2026 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
loglevel=error
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:idxinit]
command=/usr/local/bin/index-init.sh 180
autostart=true
autorestart=false
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
directory=/data
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:maps]
command=/usr/bin/goStatic
  -path /opt/maps
  -port %(ENV_OFFLINE_REGION_MAPS_PORT)s
  -append-header 'Access-Control-Allow-Origin:*'
  -enable-health
  -vhost ''
autostart=true
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
directory=/opt/maps
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:cron]
autostart=true
autorestart=true
command=/usr/local/bin/supercronic -json "%(ENV_SUPERCRONIC_CRONTAB)s"
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
